# config/.config/systemd/user/kitty-update.service
[Unit]
Description=Update Kitty Terminal

[Service]
Type=oneshot
ExecStart=%h/Scripts/kitty-update.sh

